# examples/minimal-node/apps/server/systemd/baerly.service
# systemd unit for minimal-node.
# Install:
#   sudo cp apps/server/systemd/baerly.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now baerly
# Edit `User`, `WorkingDirectory`, `EnvironmentFile` for your host.

[Unit]
Description=Baerly app — minimal-node
Documentation=https://docs.baerly.dev/
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=baerly
Group=baerly
WorkingDirectory=/opt/baerly
EnvironmentFile=/etc/baerly/minimal-node.env
ExecStart=/usr/bin/node dist/server.js
Restart=on-failure
RestartSec=5s
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=15s
LimitNOFILE=65536
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
